# roles/kubernetes-master/templates/etcd.service.j2
[Unit]
Description=etcd
Documentation=https://coreos.com/etcd/docs/latest/
After=network.target

{% set cluster_addr = cluster_url ~ (etcd_port + 1) %}
[Service]
Environment="ETCD_UNSUPPORTED_ARCH=arm64"
ExecStart=/usr/bin/etcd \
  --name master \
  --listen-client-urls {{ etcd_addr }} \
  --advertise-client-urls {{ etcd_addr }} \
  --data-dir=/var/lib/etcd \
  --initial-advertise-peer-urls {{ cluster_addr }} \
  --listen-peer-urls {{ cluster_addr }} \
  --initial-cluster-token {{ 'etcd-cluster-' ~ (range(1,100) | random) }} \
  --initial-cluster master={{ cluster_addr }} \
  --initial-cluster-state new \
  --cert-file={{ remote_cert ~ "main.crt" }} \
  --key-file={{ remote_cert ~ "main.key" }} \
  --peer-cert-file={{ remote_cert ~ "main.crt" }} \
  --peer-key-file={{ remote_cert ~ "main.key" }} \
  --trusted-ca-file={{ remote_cert ~ "ca.crt" }} \
  --peer-trusted-ca-file={{ remote_cert ~ "ca.key" }} \
  --peer-client-cert-auth \
  --client-cert-auth 
  
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target

